# heartbeat/bareclaw-heartbeat.service
[Unit]
Description=BAREclaw heartbeat — checks for pending work
After=network.target

[Service]
Type=oneshot
ExecStart=/bin/bash "__HEARTBEAT_SCRIPT__"
